# systemd/sysx-monitor.service
[Unit]
Description=sysx-monitor

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opy/sysx-monitor/run.py

[Install]
WantedBy=sysx-monitor.target
