# systemd/sysx-monitor.target
[Unit]
Description=5 Seconds Timer
StopWhenUnneeded=yes
